# ops/systemd/sealai-redis-capacity-guard.timer
[Unit]
Description=Run sealingAI Redis capacity guard every fifteen minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
AccuracySec=1min
RandomizedDelaySec=30s
Persistent=true
Unit=sealai-redis-capacity-guard.service

[Install]
WantedBy=timers.target
