# ops/systemd/sealai-listener-guard.timer
[Unit]
Description=Observe sealingAI host listeners every five minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
RandomizedDelaySec=30s
Persistent=true
Unit=sealai-listener-guard.service

[Install]
WantedBy=timers.target
