# ops/systemd/sealai-disk-guard.timer
[Unit]
Description=Run sealingAI Docker disk guard every five minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
RandomizedDelaySec=20s
Persistent=true
Unit=sealai-disk-guard.service

[Install]
WantedBy=timers.target
