# ops/systemd/sealai-redis-capacity-guard.service
[Unit]
Description=sealingAI read-only Redis capacity and namespace guard
Documentation=file:///usr/local/share/doc/sealai/redis-capacity-guard.md
After=local-fs.target

[Service]
Type=oneshot
User=sealai-redis-guard
Group=sealai-redis-guard
SupplementaryGroups=redis
UMask=0077
LoadCredentialEncrypted=redis-guard-config
LoadCredentialEncrypted=redis-auth
StateDirectory=sealai-redis-capacity-guard
StateDirectoryMode=0700
RuntimeDirectory=sealai-redis-capacity-guard
RuntimeDirectoryMode=0700
ExecStart=/usr/local/libexec/sealai/redis-capacity-guard.sh --config=%d/redis-guard-config --credential-file=%d/redis-auth check
SuccessExitStatus=10
TimeoutStartSec=10min
MemoryHigh=192M
MemoryMax=256M
TasksMax=32
NoNewPrivileges=yes
CapabilityBoundingSet=
DevicePolicy=closed
LockPersonality=yes
MemoryDenyWriteExecute=yes
PrivateDevices=yes
PrivateTmp=yes
ProtectClock=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectProc=invisible
ProcSubset=pid
ProtectSystem=strict
ReadWritePaths=/var/lib/sealai-redis-capacity-guard /run/sealai-redis-capacity-guard
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SocketBindDeny=any

[Install]
WantedBy=multi-user.target
