# seo/systemd/sealai-seo-backup.service
[Unit]
Description=SealAI SEO SQLite backup

[Service]
Type=oneshot
WorkingDirectory=/home/thorsten/sealai
ExecStart=/home/thorsten/sealai/seo/scripts/backup_sqlite.sh
