# seo/systemd/sealai-seo-backup.timer
[Unit]
Description=Weekly SealAI SEO SQLite backup timer

[Timer]
OnCalendar=Sun *-*-* 05:45:00
Persistent=true

[Install]
WantedBy=timers.target
