# seo/systemd/sealai-seo-gsc-sync.service
[Unit]
Description=SealAI SEO GSC sync

[Service]
Type=oneshot
WorkingDirectory=/home/thorsten/sealai
ExecStart=/home/thorsten/sealai/seo/scripts/run_gsc_sync.sh
