# ops/systemd/user/sealai-rag-paperless-sync.service
[Unit]
Description=SealAI Paperless to RAG sync

[Service]
Type=oneshot
ExecStart=%h/.local/bin/sealai-rag-paperless-sync
