# ops/systemd/user/sealai-rag-paperless-sync.timer
[Unit]
Description=Run SealAI Paperless to RAG sync regularly

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true
Unit=sealai-rag-paperless-sync.service

[Install]
WantedBy=timers.target
