# ops/sealai-stack.service
[Unit]
Description=Start SealAI backend and keycloak services
After=network-online.target docker.service ufw.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/home/thorsten/sealai/ops/up-prod.sh

[Install]
WantedBy=multi-user.target
